# deploy/systemd/monitor-health.service
[Unit]
Description=GitHub Repo Monitor healthcheck
After=network-online.target

[Service]
Type=oneshot
User=monitor
Group=monitor
WorkingDirectory=/opt/monitor
Environment=MONITOR_DB_PATH=/var/lib/monitor/monitor.db
EnvironmentFile=/etc/monitor/monitor.env
ExecStart=/opt/monitor/venv/bin/python /opt/monitor/scripts/healthcheck.py
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/var/lib/monitor
